# supervisor-ohmygoa.conf
[program:ohmygoa-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php8.2 /var/www/ohmygoa/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/ohmygoa/storage/logs/worker.log
stopwaitsecs=3600
user=www-data

[program:ohmygoa-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do /usr/bin/php8.2 /var/www/ohmygoa/artisan schedule:run --verbose --no-interaction & sleep 60; done"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/ohmygoa/storage/logs/scheduler.log
user=www-data
